# roles/git_mirror/templates/mirror.service.j2
[Unit]
Description=Mirror {{ NAME }} repo

[Service]
Type=oneshot
# Note the shell script must have shebang
ExecStart={{ project_path }}/sync_mirror.sh
